# cookbooks/supervisor/templates/default/worker.conf.erb
[program:<%= @name %>]
<% if @directory %>
directory = <% @directory %>
<% end %>
command = <%= @command %>
user = <%= @user %>
environment = <%= @env.map{|k, v| "#{k}=\"#{v}\""}.join(',') %>
<% if @stopsignal %>
stopsignal = <%= @stopsignal %>
stopwaitsecs = <%= @stopwaitsecs %>
<% end %>
numprocs = <%= @numprocs %>
autostart = <%= @autostart %>
autorestart = <%= @autorestart %>
process_name = <%= @name %>_%(process_num)s
stdout_logfile = <%= @log_dir %>/<%= @name %>_%(process_num)s_stdout.log
stderr_logfile = <%= @log_dir %>/<%= @name %>_%(process_num)s_stderr.log
