# cookbooks/supervisor/templates/default/worker.conf.erb
[program:<%= @name %>]
command = sudo -H -u <%= @user %> <%= @command %>
numprocs = <%= @numprocs %>
autostart = <%= @autostart %>
process_name = <%= @name %>_%(process_num)s