# cookbooks/graphite/templates/default/carbon_init_d.erb
#!/bin/bash -e

### BEGIN INIT INFO
# Provides: carbon
# Required-Start:
# Required-Stop:
# Default-Start:  2 3 4 5
# Default-Stop:
# Short-Description: carbon
### END INIT INFO

PID_FILE="<%= @graphite_directory %>/storage/carbon-cache-a.pid"
cd <%= @graphite_directory %>

<% if @whisper_dev_shm_size %>
mkdir -p /dev/shm/whisper
chown -R www-data /dev/shm/whisper
mount -o remount,size=<%= @whisper_dev_shm_size %> /dev/shm
<% end %>

#if [ "$1" = "start" ]; then
#  if [ -f "$PID_FILE" ]; then
#    if [ ! -d $(cat $PID_FILE) ]; then
#      rm $PID_FILE
#    fi
#  fi
#fi

pid=$(cat $PID_FILE)

if [[ "$(pgrep -f \"<%= @pypy %>) ./bin/carbon-cache.py start\"" == "$pid" ]]; then
  echo "running"
else
  sudo -u www-data <%= @pypy %> ./bin/carbon-cache.py $*
fi


if [ "$1" = "stop" ]; then
 count=0
 while true; do
  if [ "$count" = "60" ]; then
    echo "Unable to stop carbon killing it"
    kill -9 $pid
    break
  fi
  if [ -f "$PID_FILE" ]; then
   if [ -d "/proc/$pid" ]; then
     echo "."
     sleep 1
   else
     break
   fi
  else
   break
  fi
  count=$(($count + 1))
 done
fi
