# cookbooks/graphite/templates/default/carbon_init_d.erb
#!/bin/sh -e

### BEGIN INIT INFO
# Provides: graphite
# Required-Start:
# Required-Stop:
# Default-Start:  2 3 4 5
# Default-Stop:
# Short-Description: graphite
### END INIT INFO

cd <%= @graphite_directory %> && sudo -u www-data ./bin/carbon-cache.py $*