# cookbooks/nodejs/templates/default/init_d.erb
#! /bin/bash

### BEGIN INIT INFO
# Provides:   <%= @name %>
# Required-Start:
# Required-Stop:
# Default-Start:  2 3 4 5
# Default-Stop:
# Short-Description:  <%= @name %>
### END INIT INFO

set -e

NODE_VERSION_FILE="<%= @app_path %>/.node_version"
SCRIPT="<%= @app_path %>/<%= @script %>"
NAME="<%= @name %>"
DESC="node"
PID_FILE="<%= @pid_files_path %>/<%= @pid_file_name %>.pid"
START_STOP_DAEMON="/sbin/start-stop-daemon"

test -f $NODE_VERSION_FILE || exit 0
test -f $SCRIPT || exit 0

. <%= @user_home %>/.nvm/nvm.sh || true
nvm use `cat $NODE_VERSION_FILE`

DAEMON=`which node`

<%= @startup_code %>

test -x $DAEMON || exit 0

umask 022

. /lib/lsb/init-functions

[ -f /etc/default/<%= @name %> ] && . /etc/default/<%= @name %>

NODE_OPTS="$NODE_OPTS <% if @cluster_mode %> --cluster_pidfiles=<%= @pid_files_path %> <% end %>"

status_nodejs() {
  status_of_proc -p "$PID_FILE" "node" "$NAME"
}

kill_nodejs() {
  SIGNAL=$1

  if [ ! "$PID_FILE" = "" ]; then
    if [ -f $PID_FILE ]; then
      kill $SIGNAL `cat $PID_FILE` || true
    fi
  fi
}

start_nodejs() {
  $START_STOP_DAEMON -b -c <%= @user %> <% unless @cluster_mode %> -m <% end %> --start --pidfile $PID_FILE --exec $DAEMON -- $SCRIPT $NODE_OPTS
  <% unless @cluster_mode %> chown <%= @user %> $PID_FILE <% end %>
  sleep 1
  status_nodejs

}

case "$1" in
  start)
    echo -n "Starting $DESC: "
    start_nodejs
    echo "$NAME."
    ;;
  stop)
    echo -n "Stopping $DESC: "
    kill_nodejs
    echo "$NAME."
    ;;
  restart)
    echo -n "Restarting $DESC: "
    kill_nodejs || true
    sleep 1
    start_nodejs
    echo "$NAME."
    ;;
  reload)
    echo -n "Reloading $DESC: "
    status_nodejs || USE_START=1
    if [ "$USE_START" = "" ]; then
      kill_nodejs "-USR2"
    else
      echo "Process not found"
      echo -n "Starting $DESC: "
      start_nodejs
    fi
    echo "$NAME."
    ;;
  status)
    status_nodejs
    ;;
  *)
    echo "Usage: /etc/init.d/<%= @name %> {start|stop|restart|reload|status}"
    exit 1
esac

exit 0
