# cookbooks/rails/templates/default/resque.conf.erb
[program:resque] 
command = sudo -u <%= @user %> <%= @command %>
numprocs = <%= @numprocs %>
autostart = <%= @autostart %>
process_name = resque%(process_num)s